# deploy/private-drop-agent.service.example
[Unit]
Description=Private Drop shell agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/private-drop/private-drop-agent --config /etc/private-drop-agent/agent.toml
Restart=always
RestartSec=5

# Keep environment minimal; agent behavior should come from agent.toml.
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
