# deploy/webcodex-agent.service.example
# systemd unit for the WebCodex agent.
#
# Copy to /etc/systemd/system/webcodex-agent.service, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now webcodex-agent
#
# The agent owns local machine knowledge (project id/path/policy/capabilities)
# and registers it with the server at startup. The server is a zero-project-
# config relay and does NOT need a server-side projects.toml for runtime project
# discovery; project ids are exposed as agent:<client_id>:<project_id>.
#
# New deployments should prefer transport = "auto" in agent.toml with [quic]
# configured. The agent tries QUIC first, then WebSocket, then polling. Use a
# strict transport only when you intentionally want exactly one transport.

[Unit]
Description=WebCodex agent (auto transport: QUIC, WebSocket, polling)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/webcodex/webcodex-agent --config /etc/webcodex/agent.toml
Restart=always
RestartSec=5

# Keep environment minimal; agent behavior should come from agent.toml.
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
